# scripts/local-server/gcs-forensic-push.timer
[Unit]
Description=BwAI forensic log push to GCS (5 min)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
